# erp/deploy/systemd/artifact-erp.service
[Unit]
Description=Artifact ERP Backend — Fastify API Server
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=adam
Group=adam
WorkingDirectory=/home/adam/workspace/business_erp/backend
ExecStart=/usr/bin/npx tsx src/index.ts
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# Environment
EnvironmentFile=/home/adam/workspace/business_erp/backend/.env

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=artifact-erp

[Install]
WantedBy=multi-user.target
